# src/modules/fullpageos/filesystem/root_init/etc/systemd/system/cec-schedule.timer
[Unit]
Description=Run CEC TV schedule check every minute

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
